# root-skopeo/etc/periodic/hourly/skopeo-sync
#!/bin/bash -e

LOG=/tmp/skopeo-sync.log

s6-setuidgid abc echo $(date) === skopeo sync > ${LOG}
/app/skopeo-sync >> ${LOG} 2>&1
