# root/etc/periodic/daily/registry-gc
#!/bin/bash -e

LOG=/tmp/registry-gc.log

s6-setuidgid abc echo $(date) === registry garbage-collect > ${LOG}
/app/registry-gc >> ${LOG} 2>&1
