# root/etc/cron.daily/registry-gc
#!/bin/bash -e

s6-setuidgid abc \
    registry garbage-collect -m ${REG_CONFIG_PATH}/config.yml
