# root/etc/cron.daily/registry-gc
#!/bin/bash -e

registry garbage-collect -m ${REGISTRY_CONFIG_PATH}/config.yml
