# root/etc/s6-overlay/s6-rc.d/svc-registry/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc \
    registry serve ${REG_CONFIG_PATH}/config.yml
