# root/etc/s6-overlay/s6-rc.d/init-auth-config/run
#!/usr/bin/with-contenv bash

htpasswd -Bbn ${AUTH_USER} ${AUTH_PASS} > ${AUTH_CONFIG_PATH}/htpasswd
